# config/launchd/com.dotfiles.bootstrap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles.bootstrap</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/Shared/dotfiles/bin/dotfiles-bootstrap-agent</string>
    </array>
    <!-- Run on every login. bootstrap.zsh is idempotent: no-ops on second run. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <!-- Low priority — bootstrap is I/O light and non-urgent at login. -->
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
